# deploy/systemd/euthergate.service
[Unit]
Description=EutherGate remote development gateway
After=euthergate-forge.service euthergate-tmux.service
Wants=euthergate-forge.service euthergate-tmux.service

[Service]
Type=simple
WorkingDirectory=%h/EutherGate
ExecStart=%h/EutherGate/scripts/run-service.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
